# ssd1.mount
[Unit]
Description=Mount ssd1 disk

[Install]
WantedBy=multi-user.target

[Mount]
What=/dev/sdc1
Where=/mnt/ssd1
Type=mount
