# ssd2.mount
[Unit]
Description=Mount ssd2 disk

[Install]
WantedBy=multi-user.target

[Mount]
What=/dev/sdd1
Where=/mnt/ssd2
Type=mount
